# contrib/init/shieldd.conf
description "SHIELD Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env SHIELDD_BIN="/usr/bin/shieldd"
env SHIELDD_USER="shield"
env SHIELDD_GROUP="shield"
env SHIELDD_PIDDIR="/var/run/shieldd"
# upstart can't handle variables constructed with other variables
env SHIELDD_PIDFILE="/var/run/shieldd/shieldd.pid"
env SHIELDD_CONFIGFILE="/etc/shield/shield.conf"
env SHIELDD_DATADIR="/var/lib/shieldd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # shieldd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$SHIELDD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run shieldd."
        echo "The setting must appear in $SHIELDD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"SHIELD Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$SHIELDD_PIDDIR"
    chmod 0755 "$SHIELDD_PIDDIR"
    chown $SHIELDD_USER:$SHIELDD_GROUP "$SHIELDD_PIDDIR"
    chown $SHIELDD_USER:$SHIELDD_GROUP "$SHIELDD_CONFIGFILE"
    chmod 0660 "$SHIELDD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$SHIELDD_PIDFILE" \
    --chuid $SHIELDD_USER:$SHIELDD_GROUP \
    --exec "$SHIELDD_BIN" \
    -- \
    -pid="$SHIELDD_PIDFILE" \
    -conf="$SHIELDD_CONFIGFILE" \
    -datadir="$SHIELDD_DATADIR" \
    -disablewallet \
    -daemon

